# modules/upgrade/local.dotfiles.upgrade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.dotfiles.upgrade</string>

  <key>ProgramArguments</key>
  <array>
    <string>__HOME__/.local/bin/dotfiles-upgrade</string>
  </array>

  <!-- Weekly: Monday at 12:00. Runs at next wake if the Mac was asleep. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>12</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/dotfiles-upgrade.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/dotfiles-upgrade.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
